# share/wolf-stat-client.cron
# 

SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

* *    * * *   wolfstat  /opt/wolf-stat-client/3.0/bin/wolf-stat-client SERVERNAME >/dev/null 2>&1

